# systemd/runfast.service
[Unit]
Description=Start fast pipeline (watcher → checker → auth)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl start authLooperBackend.service
