# systemd/runslow.timer
[Unit]
Description=Slow inscription indexer timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=1800s      # 30 min (adjust as you like)
Unit=indexLooper.service   # ← tell systemd which service to fire
Persistent=true

[Install]
WantedBy=timers.target
